# deploy/systemd/adam-asr-whisperx.service
[Unit]
Description=Adam Chip · WhisperX ASR (Docker, CUDA on Jetson)
Documentation=https://github.com/m-bain/whisperX
# CRITICAL: This unit MUST run the Dockerized WhisperX container — it is the only
# path with working CUDA on aarch64. The native Python ASR (Speech.ASR_WhisperX)
# silently falls back to CPU because ctranslate2 pip-wheel has no CUDA build for
# aarch64. CPU runtime gives ~12s/transcribe → realtime UX collapses (see Test 13
# entry in tests/RESULTS.md).
#
# The compose service `adam-asr-whisperx` declares NVIDIA runtime + a healthcheck
# that asserts `/health → device=cuda`. If the container ever falls back, Docker
# will mark it unhealthy and the orchestrator will see the failure cleanly.
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
User=i17jet
Group=i17jet
WorkingDirectory=/home/i17jet/Agents/Adam-Chip

# Idempotent: `docker compose up -d` is a no-op if container is already running.
ExecStart=/usr/bin/docker compose -f /home/i17jet/Agents/Adam-Chip/compose.yaml up -d adam-asr-whisperx
ExecStop=/usr/bin/docker compose -f /home/i17jet/Agents/Adam-Chip/compose.yaml stop adam-asr-whisperx

# Container itself has restart: unless-stopped — Docker daemon handles per-process
# restarts. Systemd only orchestrates lifecycle (up at boot, stop on shutdown).
TimeoutStartSec=300
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adam-asr-whisperx

[Install]
WantedBy=multi-user.target
