# deploy/systemd/adam-tts-silero.service
[Unit]
Description=Adam Chip Silero TTS Service
After=network-online.target sound.target user@1000.service
Wants=network-online.target sound.target user@1000.service
# StartLimit must live in [Unit] — systemd silently ignores it under [Service].
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
Type=simple
User=i17jet
Group=i17jet
SupplementaryGroups=audio
WorkingDirectory=/home/i17jet/Agents/Adam-Chip
EnvironmentFile=-/etc/adam-chip/adam.env
Environment=PYTHONPATH=/home/i17jet/Agents/Adam-Chip/System
Environment=ADAM_VENV=/home/i17jet/Agents/Adam-Chip/.venv
Environment=ADAM_TTS_PORT=8082
Environment=ADAM_TTS_OUTPUT_DEVICE=plughw:1,3
# N5: 24000 = меньше WAV-байтов = быстрее transport (TFFV −100–200ms).
# Silero v5_5_ru поддерживает: 8000, 24000, 48000.
Environment=ADAM_TTS_SAMPLE_RATE=24000
Environment=ADAM_MODELS_DIR=/home/i17jet/Agents/Adam-Chip/Subsystem/Models
Environment=ADAM_TTS_MODEL_PATH=/home/i17jet/Agents/Adam-Chip/Subsystem/Models/silero/v5_5_ru.pt
# Lets the service push audio through the operator's PulseAudio session.
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
ExecStartPre=/usr/bin/test -x /home/i17jet/Agents/Adam-Chip/.venv/bin/python
# Start PulseAudio for this user session if not already running.
ExecStartPre=-/usr/bin/pulseaudio --start --daemonize=yes
ExecStart=/home/i17jet/Agents/Adam-Chip/.venv/bin/python -m Speech.TTS
Restart=on-failure
RestartSec=5
TimeoutStartSec=60

[Install]
WantedBy=multi-user.target
