# deploy/systemd/adam-asr-whisperx.service
[Unit]
Description=Adam Chip WhisperX ASR Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=i17jet
WorkingDirectory=/home/i17jet/Agents/Adam-Chip
Environment=PYTHONPATH=System
Environment=PATH=/home/i17jet/Agents/Adam-Chip/.venv/bin
Environment=ADAM_ASR_PORT=8095
Environment=ADAM_ASR_WHISPERX_MODEL=large-v3
Environment=ADAM_ASR_LANGUAGE=ru
Environment=ADAM_ASR_DEVICE=cuda
Environment=ADAM_ASR_COMPUTE_TYPE=float16
Environment=HF_HOME=/home/i17jet/Agents/Adam-Chip/data/hf_cache
ExecStart=/home/i17jet/Agents/Adam-Chip/.venv/bin/python -m Speech.ASR_WhisperX
Restart=always
RestartSec=5
TimeoutStopSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=adam-asr-whisperx

# GPU access
DeviceAllow=/dev/nvidia0 rw
DeviceAllow=/dev/nvidiactl rw

[Install]
WantedBy=multi-user.target
