# deploy/systemd/adam-llm.service
[Unit]
Description=Adam Chip · llama-server (LLM inference, OpenAI-compatible API :8081)
Documentation=https://github.com/ggerganov/llama.cpp
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=i17jet
Group=i17jet
SupplementaryGroups=render video

WorkingDirectory=/home/i17jet/Agents/Adam-Chip
EnvironmentFile=-/etc/adam-chip/adam.env

# Single source of truth for the model is System/Config.json
# (services.llm.model_path + services.llm.speculative). The launcher
# script reads it and exec's llama-server with the right flags.
# Env vars below are low-priority fallbacks used only when Config.json
# omits the field (and ADAM_LLM_GGUF_PATH is honored if model_path is
# absent — kept for backwards compatibility).
Environment=ADAM_LLM_LLAMACPP_DIR=/home/i17jet/Agents/Adam-Chip/Subsystem/llama.cpp
Environment=ADAM_LLM_PORT=8081
Environment=ADAM_LLM_GPU_LAYERS=99
Environment=ADAM_CONFIG=/home/i17jet/Agents/Adam-Chip/System/Config.json

ExecStartPre=/usr/bin/test -x ${ADAM_LLM_LLAMACPP_DIR}/build/bin/llama-server
ExecStartPre=/usr/bin/test -f ${ADAM_CONFIG}

ExecStart=/home/i17jet/Agents/Adam-Chip/scripts/adam_llm_run.sh

Restart=on-failure
StartLimitIntervalSec=120
StartLimitBurst=5
RestartSec=5
TimeoutStartSec=180
# llama-server loads the model into GPU memory on startup — can take 30–120s.

# cgroup limits: leave ~1 GB for OS and other services on 16 GB Jetson.
MemoryMax=15G
CPUQuota=90%

[Install]
WantedBy=multi-user.target
