# deploy/systemd/adam-asr-speaches.service
[Unit]
Description=Adam Chip ASR speaches (faster-whisper + CUDA, OpenAI-compatible API)
After=network-online.target docker.service
Requires=docker.service
Wants=network-online.target

[Service]
Type=simple
User=i17jet
Group=i17jet
WorkingDirectory=/home/i17jet/Agents/Adam-Chip
EnvironmentFile=-/etc/adam-chip/adam.env
Environment=ADAM_ASR_IMAGE=dustynv/speaches:r36.4.0-cu128-24.04
Environment=ADAM_ASR_MODEL=Systran/faster-whisper-base
Environment=ADAM_ASR_PORT=8083
Environment=HF_HUB_OFFLINE=1
Environment=TRANSFORMERS_OFFLINE=1
ExecStartPre=/bin/bash -c 'docker image inspect dustynv/speaches:r36.4.0-cu128-24.04 >/dev/null 2>&1 || docker pull dustynv/speaches:r36.4.0-cu128-24.04'
ExecStart=/home/i17jet/Agents/Adam-Chip/scripts/adam_asr_speaches.sh fg
ExecStop=/home/i17jet/Agents/Adam-Chip/scripts/adam_asr_speaches.sh stop
Restart=on-failure
StartLimitIntervalSec=120
StartLimitBurst=5
RestartSec=5
TimeoutStartSec=300
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
