# deploy/systemd/adam-orchestrator.service
[Unit]
Description=Adam Chip Orchestrator
After=network-online.target adam-llm.service docker.service
Wants=network-online.target adam-llm.service docker.service
# StartLimit must live in [Unit], not [Service] (systemd silently ignores it
# under [Service]). Without this, a port-collision crashloop (EADDRINUSE on
# 8080 after stale process holds the socket) was running 27+ restarts in 6 min
# during T14c, each cycle re-triggering boot-time audio side effects.
StartLimitIntervalSec=120
StartLimitBurst=3

[Service]
Type=simple
User=i17jet
Group=i17jet
SupplementaryGroups=audio video render docker
WorkingDirectory=/home/i17jet/Agents/Adam-Chip
EnvironmentFile=-/etc/adam-chip/adam.env
Environment=PYTHONPATH=/home/i17jet/Agents/Adam-Chip/System
Environment=ADAM_VENV=/home/i17jet/Agents/Adam-Chip/.venv
# Lets the service push/pull audio through the operator's PulseAudio session.
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=PULSE_SERVER=unix:/run/user/1000/pulse/native
PermissionsStartOnly=true
ExecStartPre=/home/i17jet/Agents/Adam-Chip/scripts/adam_power_maxn.sh
# Start PulseAudio for this user session if not already running.
ExecStartPre=-/usr/bin/pulseaudio --start --daemonize=yes
ExecStartPre=/usr/bin/test -x /home/i17jet/Agents/Adam-Chip/.venv/bin/python
ExecStartPre=/bin/sh -c 'until curl --noproxy "*" -fsS http://127.0.0.1:8081/health >/dev/null 2>&1; do sleep 5; done'
ExecStart=/home/i17jet/Agents/Adam-Chip/.venv/bin/python /home/i17jet/Agents/Adam-Chip/System/Orchestrator.py
Restart=on-failure
# RestartSec=10 gives the kernel time to release port 8080 from TIME_WAIT
# after a SIGTERM'd previous process. Previously 3s caused EADDRINUSE loops.
RestartSec=10
TimeoutStartSec=300
TimeoutStopSec=15
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
